# scripts/rpm/preremove
#!/bin/sh
# RPM %preun — $1 is 0 on final removal, 1 on upgrade. Only tear down on removal.
set -e
if [ "$1" = 0 ]; then
    if systemctl is-active --quiet iprd.service; then
        systemctl stop iprd.service
    fi
    systemctl disable iprd.service || true
fi
